# debian/vyatta-ravpn.postinst.in
#!/bin/bash

prefix=@prefix@
exec_prefix=@exec_prefix@
sysconfdir=@sysconfdir@
bindir=@bindir@
sbindir=@sbindir@

# remove init of daemons that we start/stop
for init in xl2tpd pptpd; do
  update-rc.d -f ${init} remove >/dev/null
done

for cfg in /etc/ipsec.d/tunnels/remote-access \
           /etc/ppp/options.xl2tpd \
           /etc/ppp/options.pptpd; do
  mkdir -p ${cfg%/*}
  touch $cfg
done

mkdir -p ${sysconfdir}/ravpn/sessions

for auth in auth-up auth-down; do
  [ -f /etc/ppp/$auth ] && [ ! -f /etc/ppp/$auth.vyatta-save ] \
    && mv -f /etc/ppp/$auth /etc/ppp/$auth.vyatta-save
  cp -f ${sysconfdir}/ravpn/ppp-$auth /etc/ppp/$auth
done

